# files/default/nat.conf
description "Configure iptables rules for NAT"
author "Amazon Linux AMI <amazon-linux-ami@amazon.com>"

start on started network
task
chdir /tmp
exec /usr/sbin/configure-pat.sh
